# install/pi/systemd/netguard-api.service
[Unit]
Description=NetGuard REST API and Dashboard
Documentation=https://github.com/netguard/netguard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=netguard
Group=netguard
WorkingDirectory=/opt/netguard
EnvironmentFile=-/etc/netguard/netguard.env
ExecStart=/opt/netguard/venv/bin/python -m uvicorn api.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=netguard.target
