# install/pi/systemd/netguard-inbound-detector.service
[Unit]
Description=NetGuard Inbound Connection Detector
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/netguard
EnvironmentFile=/etc/netguard/netguard.env
ExecStart=/opt/netguard/venv/bin/python daemon/detection/inbound_connection_detector.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=netguard.target
