# install/pi/systemd/netguard-dns-monitor.service
[Unit]
Description=NetGuard DNS Monitor (Scapy + dnsmasq log)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/netguard
EnvironmentFile=/etc/netguard/netguard.env
ExecStart=/opt/netguard/venv/bin/python daemon/scanner/dns_monitor.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=netguard.target
