# install/pi/systemd/netguard-msp-agent.service
[Unit]
Description=NetGuard MSP agent heartbeat
After=network-online.target

[Service]
Type=oneshot
User=netguard
WorkingDirectory=/opt/netguard
Environment=PYTHONPATH=/opt/netguard/daemon
EnvironmentFile=-/etc/netguard/netguard.env
ExecStart=/opt/netguard/venv/bin/python daemon/msp/agent_client.py

[Install]
WantedBy=multi-user.target
